# templates/default/init.logstash.conf.erb
# logstash - agent instance
#
# This file is based on https://github.com/elasticsearch/logstash/blob/master/pkg/logstash.upstart.ubuntu
# It is developed for upstart v0.6.5 (an old version) because CentOS ships with an old upstart version by default.

description     "logstash agent"

start on virtual-filesystems
stop on runlevel [06]

# Respawn it if the process exits
respawn

# We're setting high here, we'll re-limit below.
limit nofile 65550 65550

# You need to chdir somewhere writable because logstash needs to unpack a few
# temporary files on startup.
chdir <%= node['logstash']['homedir'] %>
script
    # Defaults
    PATH="/bin:/usr/bin"
    LS_HOME="<%= node['logstash']['homedir'] %>"
    LS_HEAP_SIZE="500m"
    LS_JAVA_OPTS="-Djava.io.tmpdir=${LS_HOME}"
    LS_USE_GC_LOGGING=""
    LS_OPEN_FILES=16384

    # Override our defaults with user defaults:
    [ -f /etc/default/logstash ] && . /etc/default/logstash

    HOME="${HOME:-$LS_HOME}"
    JAVA_OPTS="${LS_JAVA_OPTS}"

    # Export variables
    export PATH HOME JAVA_OPTS LS_HEAP_SIZE LS_JAVA_OPTS LS_USE_GC_LOGGING

    sudo -u <%= node['logstash']['user'] %> -s /bin/sh <%= node['logstash']['dir'] %>/bin/logstash agent -f "<%= node['logstash']['conf_dir'] %>" -l "<%= node['logstash']['log_file'] %>"
    emit logstash-running
end script